# v2socks.service
[Unit]
Description=v2socks

[Service]
Type=simple
ExecStart=/usr/bin/v2socks plain